# systemd/arnold-random-quotes-stop.service
# arnold-random-quotes stop systemd timer unit file

[Unit]
Description=Stop arnold-random-quotes every day at 00:00
SuccessAction=exit
